# hubitat-go-mcp-sse.service
[Unit]
Description=Hubitat Go MCP Server (SSE Mode)
After=network.target
Documentation=https://github.com/alexj212/hubitat-go-mcp

[Service]
Type=simple
User=alexj
Group=alexj
WorkingDirectory=/docker/hubitat-go-mcp
ExecStart=/docker/hubitat-go-mcp/hubitat-go-mcp -mode=sse
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hubitat-go-mcp-sse

# Environment file
EnvironmentFile=/docker/hubitat-go-mcp/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/docker/hubitat-go-mcp

[Install]
WantedBy=multi-user.target
